# systemd/monk.service
[Unit]
Description=Monk daemon
After=network.target

[Service]
ExecStart=/usr/bin/monk daemon
Restart=always

[Install]
WantedBy=basic.target
